# src/n3rv/init/templates/n3rv/systemd/n3rv-hub.service.j2
[Unit]
Description=NERV A2A Hub for {{ project_name }}
After=network.target

[Service]
Type=simple
WorkingDirectory={{ project_root }}
ExecStart={{ n3rv_binary }} hub start
Restart=on-failure
RestartSec=5
StandardOutput=append:{{ project_root }}/.n3rv/logs/hub.log
StandardError=append:{{ project_root }}/.n3rv/logs/hub.log

[Install]
WantedBy=default.target
